# sockperf-server.service
[Unit]
Description=sockperf server
After=network.target

[Service]
ExecStart=/home/adam/testing-server/sockperf server -p 20001
RestartSec=3
Restart=always

[Install]
WantedBy=multi-user.target
